# icinga2-master/rootfs/etc/s6/icinga2-cert-service/run
#! /bin/sh

PROGARGS=""
PROGHOME="/"
PROGNAME="/usr/local/bin/rest-service.rb"

exec ${PROGNAME} ${PROGARGS}
