# rootfs/etc/supervisor/conf.d/icinga2.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
loglevel=info
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:icinga2]
autostart=true
autorestart=true
command=/bin/bash -c "/opt/supervisor/icinga2_supervisor"
redirect_stderr=true
stdout_logfile=/tmp/icinga2_supervisord.log
