# rootfs/etc/supervisor.d/icinga2-cert-service.ini

[program:icinga2-cert-service]
command        = /usr/local/bin/rest-service.rb
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true

